# docker/examplewebapp/mongo.sv.conf
[program:mongo]
command=service mongodb start
