# docker/ssh/ssh.sv.conf
[program:ssh]
command=/bin/bash -c "/usr/sbin/sshd -D"
